# deploy/helia-connectivity-lab.service
[Unit]
Description=helia-connectivity-lab libp2p echo + circuit relay
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/helia-connectivity-lab
Environment=NODE_ENV=production
Environment=RELAY_LISTEN_IPV4=0.0.0.0
Environment=RELAY_DISABLE_IPV6=true
# Ports aligned with Nym exit policy (TCP): *:80-81, *:8008, *:8080, *:5000-5005; UDP QUIC 5000; UDP WebRTC 3478-3484
Environment=RELAY_TCP_PORT=81
Environment=RELAY_WS_PORT=8080
Environment=RELAY_QUIC_PORT=5000
Environment=RELAY_WEBRTC_PORT=3478
Environment=RELAY_KEY_FILE=/var/lib/helia-connectivity-lab/peer.key
Environment=RELAY_CONTROL_HTTP_PORT=8008
# Optional: GET /ipfs/<cid> on the same HTTP port as control (/health, /status).
# Environment=RELAY_IPFS_GATEWAY=1
# Environment=RELAY_IPFS_GATEWAY_LOG=1
# Environment=LIBP2P_CONN_LOG=1
# Standalone gateway only if control HTTP is disabled (port must not be RELAY_WS_PORT):
# Environment=RELAY_IPFS_HTTP_PORT=8444
# Environment=RELAY_IPFS_TLS_CERT=/path/to/fullchain.pem
# Environment=RELAY_IPFS_TLS_KEY=/path/to/privkey.pem
# Token: set RELAY_CONTROL_TOKEN only in /etc/default/helia-connectivity-lab — systemd applies EnvironmentFile after Environment= and overrides the same key.
EnvironmentFile=-/etc/default/helia-connectivity-lab
ExecStart=/usr/bin/node dist/server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
